# gpu-exporter/shopai-gpu-exporter.service
[Unit]
Description=ShopAI GPU exporter (nvidia-smi -> Prometheus :9401)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=saiteku
# Adjust if you relocate the project. The venv is created by install_systemd_service.sh.
WorkingDirectory=/home/saiteku/Project/shopai-dashboard/gpu-exporter
Environment=GPU_EXPORTER_PORT=9401
Environment=GPU_EXPORTER_INTERVAL=2
ExecStart=/home/saiteku/Project/shopai-dashboard/gpu-exporter/.venv/bin/python /home/saiteku/Project/shopai-dashboard/gpu-exporter/gpu_exporter.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
